# files/root/etc/rc.local
#!/bin/sh

## hyphop ##

#= custom rc.local

export PATH="/opt:/tmp/opt:/usr/sbin:/usr/bin:/sbin:/bin"

echo "[i] $0 started">&2

for s in /etc/rc.local.d/*.sh; do
    [ -s $s ] && {
	logger -t "$0" "start $s"
	sh $s
    }
done

logger -t "$0" DONE

## USER custom scripts there

