# files/root/etc/init.d/fan
#!/bin/sh /etc/rc.common

## hyphop ##

START=94

USE_PROCD=1
PROG=/opt/fan

start_service() {

	procd_open_instance 
	procd_set_param command $PROG monitor
	procd_close_instance

}
