# files/root/etc/init.d/ifmonitor
#!/bin/sh /etc/rc.common

## hyphop ##

START=94

USE_PROCD=1
DEV=eth0
PROG=/opt/ifmonitor

start_service() {

	procd_open_instance 
	procd_set_param command $PROG $DEV
	procd_close_instance

}
