# files/root/etc/profile.d/01-splash.mod.sh
#!/bin/sh

## hyphop ##

cat /etc/banner.mod
