# files/root/etc/profile.d/10-custom.sh
#!/bin/sh

## hyphop ##

alias opkf='opkg --force-space --force-checksum --force-depends'
alias opkd='opkg --force-space --nodeps'

export TERM=xterm
